# job.cron
# every day 4:00 am
0 4 * * * /home/hebi/github/hn-top/job.sh

